# deploy/brainsession-pwa-https.service
[Unit]
Description=BrainSession PWA HTTPS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/myagent/.ductor/workspace/external/Brainsession_pwa
ExecStart=/run/current-system/sw/bin/bash /home/myagent/.ductor/workspace/external/Brainsession_pwa/scripts/start-local-https.sh
Restart=always
RestartSec=5
Environment=HOST=0.0.0.0
Environment=BRAINSESSION_HTTPS_IPS=100.77.57.8
Environment=PATH=/run/current-system/sw/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
